# aur/pkg/screenrotator-modified-git/etc/xdg/autostart/screenrotator-autostart.desktop
[Desktop Entry]
Type=Application
Name=Screen Rotator
Exec=/usr/bin/screenrotator
NoDisplay=true
X-KDE-StartupNotify=false
